# myfiles/tmp/crontab
* * * * * /usr/local/bin/upload_s3.sh >/var/tmp/push_localfiles.txt 2>&1
* * * * * /usr/local/bin/create_user_accounts.sh >/var/tmp/create_user_accounts.txt 2>&1

